# crontab
# m h    d m w    cmd

10 10   * * *   ~/agents/testmail.py
